# infra/launchd/com.jasonchi.pa.reminder-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jasonchi.pa.reminder-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3.12</string>
        <string>-m</string>
        <string>kernel.proactive</string>
        <string>--task</string>
        <string>check-reminders</string>
        <string>--config</string>
        <string>configs/default.yaml</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>__TELEGRAM_BOT_TOKEN__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/infra/launchd/reminder-check.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/infra/launchd/reminder-check.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
